# share/musso.service
[Unit]
Description=Musso watcher mode

[Service]
Type=simple
ExecStart=/usr/bin/env musso watch

[Install]
WantedBy=default.target
